# dist/keepitup-update-nodes.timer
[Unit]
Description=The KeepItUp Update Nodes every minute

[Timer]
OnUnitActiveSec=1m

[Install]
WantedBy=timers.target