# dist/keepitup-update-nodes.service
[Unit]
Description=The KeepItUp Update Nodes
After=network.target
PartOf=keepitup.target

[Service]
Type=oneshot
User=%SYSTEM_USER%
WorkingDirectory=%DIR%
ExecStart=%DIR%/venv/bin/python %DIR%/update_nodes.py
